# rtb-bid-engine.service
[Unit]
Description=RTB Bid Engine Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/RTB/gpt_test_3
ExecStart=/root/RTB/gpt_test_3/cmd/bid-engine/bid-engine
Restart=always
RestartSec=1
StartLimitInterval=0

# Environment файл из проекта
EnvironmentFile=/root/RTB/gpt_test_3/cmd/bid-engine/bid-engine.env

# Logging
StandardOutput=append:/var/log/rtb/bid-engine.log
StandardError=append:/var/log/rtb/bid-engine.error.log
SyslogIdentifier=rtb-bid-engine

# Security
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
